# scripts/chatbot-backup-monitor.timer
[Unit]
Description=GPT Chatbot Backup Monitoring Timer

[Timer]
# Run every hour
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
